# os updates/updates-panel/hyperwebster-update-check.service
[Unit]
Description=HyperWebster update status check (writes update-status.json)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyperwebster-update-check
