# os updates/theme-polish/hyperwebster-sddm-sync.service
[Unit]
Description=Sync SDDM greeter to Caelestia scheme change
Documentation=file:///usr/local/share/hyperwebster/theme-polish/README.md

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo -n /usr/local/bin/sddm-theme-sync
